# diboards/config/supervisord.conf
[supervisord]
logfile = /usr/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /usr/log/supervisord.pid
nodeamon = False
user = root

[program:gunicorn]
command=gunicorn --config ./config/gunicorn.conf.py diboards:app
directory=/usr/src/app/
redirect_stdout=true
redirect_stderr=true
;stdout_logfile=/usr/log
;stderr_logfile=/usr/log
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
directory=/usr/src/app/
redirect_stdout=true
redirect_stderr=true
;stdout_logfile=/usr/log
;stderr_logfile=/usr/log
autostart=true
autorestart=true